# etc/supervisor/conf/controller.conf
[program:controller]
user=root
directory=%(ENV_TFW_EHMAIN_DIR)s
command=python3 controller.py
autorestart=true